# templates/etc/systemd/system/docker.service.d/service-overrides-suse.conf.erb
[Service]
EnvironmentFile=-/etc/sysconfig/docker
<% if @daemon_environment_files %><% @daemon_environment_files.each do |param| %>EnvironmentFile=-<%= param %>
<% end %><% end -%>
ExecStart=
ExecStart=/usr/bin/<%= @docker_start_command %> $OPTIONS $DOCKER_OPTS
